# src/debianpkg/postinst
#!/bin/sh
# postinst script for nsia
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package


getPassword () {
	stty -echo
	echo "Enter the password for $username:"
	read password

	echo "Enter the password again:"
	read passwordConfirm
	stty echo
}


# 1 -- Configure NSIA
if [ -d /opt/nsia/var/database ]
then
    echo "Database is already configured"
else
	# 1.1 -- Get the username
	echo "Enter a username for the administrator account. Note that the username can contain letters, numbers, dashes and periods:"
	while :
	do
		read username
	
		if [ "${#username}" = "0" ]
		then
			isInvalid=1
		else
			isInvalid=$(echo $username | sed 's/[a-zA-Z0-9_.-]*//')
		fi
	
		if [ "${#isInvalid}" != "0" ]
		then
			echo "The username is invalid. The username must only contain letters, numbers, dashes and periods."
		elif [ "${#username}" = "0" ]
		then
			echo "The username is invalid. The username must only contain letters, numbers, dashes and periods."
		else
			break
		fi
	
		echo "Please try again:"
	
	
	done

	# 1.2 -- Get the password
	while :
	do
		getPassword
		if [ "$password" = "$passwordConfirm" ]
		then
	
			if [ "${#password}" -gt 7 ]
			then
				break
			else
				echo "The passwords must be at least 8 characters, please try again"
			fi
		
		else
			echo "The passwords did not match, please try again"
		fi
	done

	# 1.3 -- Initialize NSIA
	cd /opt/nsia/bin
	java -jar nsia.jar --install $username $username $password
fi

# 2 -- Set NSIA up to start on boot
update-rc.d nsia defaults 2>/dev/null
 
# 3 -- Start NSIA
if [ -f /etc/init.d/nsia ]
then
	chmod +x /etc/init.d/nsia
	/etc/init.d/nsia start
	echo "NSIA is configured and running on port 8080 (http://localhost:8080)"
fi

echo "Go to http://threatfactor.com to download signatures"

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
