# trunk/src/debianpkg/postinst
#!/bin/sh
# postinst script for nsia
#
# see: dh_installdeb(1)

set -e
. /usr/share/debconf/confmodule

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package

# 1 -- Configure NSIA
if [ -d /opt/nsia/var/database ]
then
    echo "Database is already configured"
    
    #Upgrade the database schema...
    echo "Upgrading..."
    cd /opt/nsia/bin
    java -jar nsia.jar --upgrade
else
	
	# 1.1 -- Get the username
	db_fset nsia/username seen false
	
	while :
	do
	
		db_input high nsia/username || true
		db_go
		
		db_get nsia/username
		username="$RET"
		
		if [ "${#username}" = "0" ]
		then
			isInvalid=1
		else
			isInvalid=$(echo $username | sed 's/[a-zA-Z0-9_.-]*//')
		fi
	
		# Make sure that username is not invalid
		if [ "${#isInvalid}" = "0" ]
		then
			# Make sure the username is not empty
			if [ "${#username}" != "0" ]
			then
				break
			fi
		fi
		
	done
	
	# 1.2 -- Get the password
	db_subst nsia/password "desc" ""
	
	while :
	do
		# Get the password
		db_fset nsia/password seen false
		db_set nsia/password ""
		db_input high nsia/password || true
		db_go
		
		db_get nsia/password
		password="$RET"
		
		# Make sure the password is long enough
		if [ "${#password}" -gt 7 ]
		then
			
			# Verify the password
			db_fset nsia/password_verify seen false
			db_set nsia/password_verify ""
			db_input high nsia/password_verify || true
			db_go
			
			db_get nsia/password_verify
			passwordConfirm="$RET"
			
			# Make sure the passwords match
			if [ "$password" = "$passwordConfirm" ]
			then
				break
			fi
			
			db_subst nsia/password "desc" "Note: Your passwords did not match, please try again"
		else
			db_subst nsia/password "desc" "Note: Your password was not long enough"
		fi
		
	done
	
	# Clear the passwords from the debconf database
	db_clear
	db_unregister nsia/password
	db_unregister nsia/password_verify
	db_purge
	
	# 1.3 -- Initialize NSIA
	cd /opt/nsia/bin
	java -jar nsia.jar --install $username $username $password
fi

# 2 -- Set NSIA up to start on boot
update-rc.d nsia defaults 2>/dev/null
 
# 3 -- Start NSIA
if [ -f /etc/init.d/nsia ]
then
	chmod +x /etc/init.d/nsia
	/etc/init.d/nsia start
fi

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
